# prod/systemd/agentag-worker.service
[Unit]
Description=AgentTag async worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/agentag/app
Environment=APP_ENV=prod
Environment=APP_DEBUG=0
ExecStart=/usr/bin/php8.4 /srv/agentag/app/bin/console messenger:consume async --time-limit=3600 --memory-limit=512M
Restart=always
RestartSec=5
User=root
Group=root
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
